# systemd/airwatch-notify.service
[Unit]
Description=AirWatch Notification Service
After=network.target

[Service]
User=raspberry
WorkingDirectory=/home/raspberry/airwatch
EnvironmentFile=/home/raspberry/airwatch/.env
ExecStart=/home/raspberry/airwatch/venv/bin/python services/notification/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
