# systemd/airwatch-alerts.service
[Unit]
Description=AirWatch Alert Engine
After=network.target

[Service]
User=raspberry
WorkingDirectory=/home/raspberry/airwatch
EnvironmentFile=/home/raspberry/airwatch/.env
ExecStart=/home/raspberry/airwatch/venv/bin/python services/alert_engine/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
